# test/terraform/kubeadm/nodes/resources/system/etc/systemd/system/kubelet.service.d/bootstrap.conf
# /etc/systemd/system/kubelet.service.d/bootstrap.conf
# (!!! To be templated by Terraform !!!)
[Unit]
# Wait for cloud-init to be done with the configuration
ConditionPathExists=/etc/kubernetes/kubelet/enable

[Service]
# REF: https://kubernetes.io/docs/reference/command-line-tools-reference/kubelet/
%{ if set_node_ip ~}
EnvironmentFile=-/run/kubelet.env
ExecStartPre=/bin/sh -c 'netplan ip leases eth1 | sed -n "s|^ADDRESS=|PRIVNET_IP=|p" > /run/kubelet.env'
%{ endif ~}
ExecStart=
ExecStart=/usr/bin/kubelet \
  --cloud-provider=external \
%{ if set_node_ip ~}
  --node-ip=$${PRIVNET_IP} \
%{ endif ~}
  --bootstrap-kubeconfig=/etc/kubernetes/kubelet/bootstrap.kubeconfig \
  --kubeconfig=/etc/kubernetes/kubelet/kubeconfig \
  --config=/etc/kubernetes/kubelet/config.yaml \
  --container-runtime-endpoint=unix:///run/containerd/containerd.sock
